# packaging/macos/com.nellie-rs.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nellie-rs.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/nellie</string>
        <string>serve</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8766</string>
        <string>--data-dir</string>
        <string>/var/lib/nellie-rs</string>
        <string>--watch</string>
        <string>/Volumes/mike-github/github,/Volumes/lono-github/github,/Volumes/mmn-github/github</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>WorkingDirectory</key>
    <string>/var/lib/nellie-rs</string>

    <key>StandardOutPath</key>
    <string>/var/lib/nellie-rs/logs/nellie.out.log</string>

    <key>StandardErrorPath</key>
    <string>/var/lib/nellie-rs/logs/nellie.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NELLIE_LOG_LEVEL</key>
        <string>info</string>
        <key>RUST_BACKTRACE</key>
        <string>1</string>
    </dict>

    <key>ProcessType</key>
    <string>Standard</string>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>

    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
